# systemd/aidb-backup-qdrant.service
[Unit]
Description=AIDB Qdrant Backup Service
Documentation=https://qdrant.tech/documentation/guides/backup/
After=network-online.target qdrant.service local-ai-qdrant.service
Wants=network-online.target

[Service]
Type=oneshot
User=%u
Group=%u
WorkingDirectory=%h/Documents/try/NixOS-Dev-Quick-Deploy

# Run the backup script
ExecStart=%h/Documents/try/NixOS-Dev-Quick-Deploy/scripts/backup-qdrant.sh

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aidb-backup-qdrant

# Security hardening
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/backups

# Resource limits
CPUQuota=50%
MemoryMax=1G
TasksMax=10

[Install]
WantedBy=multi-user.target
