# templates/systemd/ai-stack-runtime-recovery.service
# Systemd service to recover Podman runtime state on boot
# Part of: NixOS-Dev-Quick-Deploy
# LEGACY: Podman runtime recovery (deprecated). K3s is the supported path.
# Install: sudo cp this file to /etc/systemd/system/
# Enable: sudo systemctl enable ai-stack-runtime-recovery.service

[Unit]
Description=Recover Podman runtime state for AI stack
After=local-fs.target
Before=podman.service

[Service]
Type=oneshot
User=root
# Update these paths to match your installation (replace @PROJECT_ROOT@)
Environment="PODMAN_RUNTIME_UID=1000"
ExecStart=@PROJECT_ROOT@/scripts/reset-podman-runtime.sh
StandardOutput=journal
StandardError=journal
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
