# ai-stack/systemd/letsencrypt-renewal.service
[Unit]
Description=AI Stack TLS Renewal (cert-manager)
After=network.target k3s.service
Requires=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=@PROJECT_ROOT@
ExecStart=@PROJECT_ROOT@/scripts/renew-tls-certificate.sh --namespace ai-stack --k8s-dir ai-stack/kubernetes/tls
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log

[Install]
WantedBy=multi-user.target
